# supervisord.conf
[program:openshare]
command = /home/dotcloud/current/dotcloud.sh
